# infra/etc/systemd/system/bootboy.service
[Unit]
Description=BootBoy boot-time configuration
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/home/wombat/Documents/github/mellow-wombat/bin/bootboy --admin-dir /var/wombat/admin --status-file /var/wombat/admin/bootboy.status.json
RemainAfterExit=yes

# Hardening note:
# Some environments (notably certain Raspberry Pi/kernel setups) can fail with
# "failed to setup mount namespacing" when mount-namespace-based hardening
# directives are enabled (e.g., ProtectSystem/ProtectHome/PrivateTmp/ReadWritePaths).
# Keep this unit compatible and add hardening only once the target supports it.
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
